# mac-mini/relay/com.seongchun.mac-mini-relay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.seongchun.mac-mini-relay</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>set -a; source "__ENV_FILE__"; set +a; exec "__NODE__" "__SCRIPT__"</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>__LOG_DIR__/mac-mini-relay.out.log</string>
  <key>StandardErrorPath</key><string>__LOG_DIR__/mac-mini-relay.err.log</string>
  <key>ThrottleInterval</key><integer>10</integer>
</dict>
</plist>
